# udp-vpn-client-0.5/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable udp-vpn-client
systemctl start udp-vpn-client
