# udp-vpn-client-0.5/DEBIAN/prerm
#!/bin/bash
systemctl stop vpn
systemctl disable vpn
