# udp-vpn-server-0.5/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable udp-vpn-server
systemctl start udp-vpn-server
